# files/shared/etc/xdg/autostart/com.bitwarden.desktop.desktop
[Desktop Entry]
Name=Bitwarden
Exec=/usr/bin/flatpak run com.bitwarden.desktop --silent
Terminal=false
Type=Application
