# files/shared/usr/lib/systemd/user/chezmoi-apply.service
[Unit]
Description=Apply chezmoi dotfiles (first login)
ConditionPathExists=!%h/.config/chezmoi/chezmoi.toml

[Service]
Type=oneshot
ExecStart=/usr/libexec/chezmoi-init
RemainAfterExit=yes

[Install]
WantedBy=default.target
