# files/system/etc/profile.d/helm-env.sh
export HELM_PLUGINS="/usr/lib/helm/plugins:${HOME}/.local/share/helm/plugins"
