# files/gaming/etc/xdg/autostart/com.discordapp.Discord.desktop
[Desktop Entry]
Name=Discord
Exec=/usr/bin/flatpak run com.discordapp.Discord --start-minimized
Terminal=false
Type=Application
